# scim.service
[Unit]
Description=Start the input service framework panel
After=e17.service

[Service]
ExecStart=/usr/bin/scim
OOMScoreAdjust=-1000
Nice=5
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
